# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/app/supervisord.pid

[program:node_backend]
command=node /app/sandbox/node/app.js
directory=/app/sandbox/node
autostart=true
autorestart=true
priority=10
stdout_logfile=/app/logs/node_stdout.log
stderr_logfile=/app/logs/node_stderr.log
startretries=3

[program:nginx]
command=nginx -c /app/config/nginx.conf
autostart=true
autorestart=true
priority=20
stdout_logfile=/app/logs/nginx_stdout.log
stderr_logfile=/app/logs/nginx_stderr.log
startretries=3

[program:openenv_api]
command=/app/.venv/bin/uvicorn server.app:app --host 0.0.0.0 --port 8000
directory=/app/env
autostart=true
autorestart=true
priority=30
environment=PYTHONPATH="/app/env",PATH="/app/.venv/bin:%(ENV_PATH)s"
stdout_logfile=/app/logs/uvicorn_stdout.log
stderr_logfile=/app/logs/uvicorn_stderr.log
startretries=3
